# extra-features/huion/80-huion.rules
ACTION=="add" \
, ATTRS{idProduct}=="006e" \
, ATTRS{idVendor}=="256c" \
, ENV{DISPLAY}=":0" \
, RUN+="/usr/local/bin/adjust-huion.sh"
