# extra-features/user-services/services/redshift.service
[Unit]
Description=Redshift

[Service]
ExecStart=redshift

Type=simple

RestartSec=3
Restart=always

[Install]
WantedBy=default.target
