# extra-features/user-services/services/xbindkeys.service
[Unit]
Description=Xbindkeys

[Service]
ExecStart=xbindkeys -n

Type=simple
RestartSec=3
Restart=on-failure

[Install]
WantedBy=default.target
