# optional-features/user-services/services/startup-application@.service
[Unit]
Description=Autostart %i on login

[Service]
ExecStart=%i

Type=simple

RestartSec=3
Restart=on-failure

[Install]
WantedBy=default.target
