# optional-features/user-services/services/update-notify.service
[Unit]
Description=Package update notification

[Service]
ExecStart=%h/.scripts/dunst/update-notification
Type=simple

[Install]
WantedBy=default.target
