# extra-features/physlock/physlock.service
[Unit]
Description=Lock session using physlock

[Service]
ExecStart=/usr/bin/physlock -l
Type=simple
RestartSec=5

[Install]
WantedBy=sleep.target 
